# supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=php -S 0.0.0.0:80 -t /var/www/html/
autostart=true
autorestart=true
stderr_logfile=/var/log/php.err.log
stdout_logfile=/var/log/php.out.log

[program:node]
command=npm start
directory=/var/www/html/
autostart=true
autorestart=true
stderr_logfile=/var/log/node.err.log
stdout_logfile=/var/log/node.out.log
